# crontab
0,15,30,45 * * * * sudo python /var/www/cocaine.ninja/mysql_expire.py
